# configs/crontab/crontab.config
0 7,19 * * * bash ~/automation/configs/crontab/crontab.sh follow >> ~/automation/logs/cron.log 2>&1
