# configs/crontab/crontab.config
0 18 * * * bash ~/automation/crontab/crontab.sh tweet >> ~/automation/logs/cron.log 2>&1
0 19 * * * bash ~/automation/crontab/crontab.sh follow >> ~/automation/logs/cron.log 2>&1
